# gpt-sovits-tts.service
# ========================================
# GPT-SoVITS TTS Server systemd 서비스
# ========================================
#
# 설치 방법:
#   1. 이 파일을 복사:
#      sudo cp gpt-sovits-tts.service /etc/systemd/system/
#
#   2. 서비스 등록 및 시작:
#      sudo systemctl daemon-reload
#      sudo systemctl enable gpt-sovits-tts
#      sudo systemctl start gpt-sovits-tts
#
#   3. 상태 확인:
#      sudo systemctl status gpt-sovits-tts
#
#   4. 로그 확인:
#      sudo journalctl -u gpt-sovits-tts -f
#
# ========================================

[Unit]
Description=GPT-SoVITS TTS Server
After=network.target

[Service]
# ⚠️ 아래 경로를 실제 환경에 맞게 수정하세요!
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/gpt-sovits
Environment="PATH=/home/ubuntu/gpt-sovits/venv/bin"

# 서버 실행 명령어 (--mode: ws, wss, both)
ExecStart=/home/ubuntu/gpt-sovits/venv/bin/python tts_api_server.py --mode ws

# 재시작 설정
Restart=always
RestartSec=10

# 로그 설정
StandardOutput=journal
StandardError=journal

# GPU 관련 환경변수
Environment="CUDA_VISIBLE_DEVICES=0"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
